# example/example.ini
[program:example]
command=/services/example
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/logs/example.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
