# example/example.ini
[program:hellogo]
command=/services/hellogo
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/logs/hellogo.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
